# coreos/jobify-server.service
[Unit]
Description= Jobify server
Requires=jobify-mongo.service
After=jobify-mongo.service

[Service]
ExecStartPre=-/usr/bin/docker stop server
ExecStartPre=-/usr/bin/docker rm server
ExecStart=/usr/bin/docker run --name server -p 3000:3000 --link mongodb:jobify jobify-server
ExecStop=-/usr/bin/docker stop server
ExecStopPost=-/usr/bin/docker rm server
